# config/launchd/com.serokell.reindex.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Daily Reindex Job
         Runs once at 7:00 AM to rebuild context graph database
         Runs BEFORE morning brief (8am) so brief has fresh data -->

    <key>Label</key>
    <string>com.serokell.reindex</string>

    <key>ProgramArguments</key>
    <array>
        <string>__SEROKELL_ROOT__/scripts/reindex.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__SEROKELL_ROOT__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.bun/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <!-- Run at 7:00 AM every day (1 hour before brief) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Logging (stored in vault) -->
    <key>StandardOutPath</key>
    <string>__VAULT_LOGS__/reindex.out.log</string>

    <key>StandardErrorPath</key>
    <string>__VAULT_LOGS__/reindex.err.log</string>

    <!-- Don't run at load - only on schedule -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
